# gateway/acronis.conf
[supervisord]
nodaemon=true

[program:cron]
command = /usr/sbin/crond -n
process_name = %(program_name)s
user=root
directory=/

autostart=true
autorestart=true
startsecs = 5
stopwaitsecs = 3600
stopasgroup = false
killasgroup = true

stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log

[program:acronis-storage-mds]
priority=1
command = /usr/sbin/acronis-storage-mds --config-file=/etc/Acronis/acronis-storage-mds.xml

[program:acronis-storage-gateway]
priority=999
command = /usr/sbin/acronis-storage-gateway --config-file=/etc/Acronis/acronis-storage-gateway.xml

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
startretries=3
